# scripts/rc.local
#!/bin/sh -e
sudo bash /home/pi/waifu-badge/scripts/convert.sh
sudo python3.5 /home/pi/waifu-badge/mcp.py &
exit 0

